# sandsifterOS/airootfs/etc/profile.d/cowsize.sh
#!/bin/bash

sudo -S mount -o remount,size=3G /run/archiso/cowspace
export TERM=xterm-256color
